# .config/systemd/user/power-profile.service
[Unit]
Description="Set power profile based on running processes"

[Service]
ExecStart=/home/mark/bin/set-power-profile.sh
